# firefly-iii/rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

echo "Setting up Firefly III cron jobs..."

# Create crontab file
cat > /etc/cron.d/firefly-cron << EOF
# Firefly III cron jobs
# Feel free to modify the schedule to your liking
# Format: minute hour day month day_of_week command

# Run Firefly III cron command daily at 3:00 AM
0 3 * * * www-data cd /var/www/html && php artisan firefly-iii:cron > /proc/1/fd/1 2>/proc/1/fd/2

# End of crontab
EOF

chmod 0644 /etc/cron.d/firefly-cron

echo "Starting cron service..."
exec cron -f
