# rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "Nginx stopped, restarting..."
sleep 1
