# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Firefly III
# Runs the Nginx daemon
# ==============================================================================

# Wait until PHP-FPM is available
bashio::net.wait_for 9000

# Ensure proper permissions for temp directories without changing ownership
chmod -R 777 /tmp/nginx || true

# Start Nginx
bashio::log.info "Starting Nginx..."
exec nginx -g "daemon off;"
