# rootfs/etc/services.d/mysql/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "MySQL stopped, restarting..."
sleep 1
