# rootfs/etc/services.d/mysql/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting MySQL..."
exec mysqld_safe
