# addon/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash

echo "Starting Firefly III application..."

# The firefly-iii image uses Apache
if [ -x "/entrypoint.sh" ]; then
    exec /entrypoint.sh
elif [ -x "/init.sh" ]; then
    exec /init.sh
else
    # Fallback to starting Apache directly
    exec apache2-foreground
fi
