# rootfs/etc/services.d/php-fpm/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "PHP-FPM stopped, restarting..."
sleep 1
