# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start PHP-FPM service
# ==============================================================================

# Display PHP information for debugging
echo "Available PHP binaries:"
find /usr -name "php*" -type f 2>/dev/null || echo "No PHP binaries found"
echo "PHP version: $(php -v 2>/dev/null || echo 'PHP not found')"

# Ensure permissions on critical directories for PHP
find /var/www/html -type d -exec chmod 777 {} \; 2>/dev/null || true
find /var/www/html -type f -exec chmod 666 {} \; 2>/dev/null || true

# Try each possible PHP-FPM executable in order without -R flag
if [ -f /usr/sbin/php-fpm8 ]; then
    echo "Using /usr/sbin/php-fpm8"
    exec /usr/sbin/php-fpm8 --nodaemonize --fpm-config /etc/php8/php-fpm.conf
elif [ -f /usr/sbin/php-fpm7 ]; then
    echo "Using /usr/sbin/php-fpm7"
    exec /usr/sbin/php-fpm7 --nodaemonize --fpm-config /etc/php7/php-fpm.conf
elif [ -f /usr/bin/php-fpm8 ]; then
    echo "Using /usr/bin/php-fpm8"
    exec /usr/bin/php-fpm8 --nodaemonize
elif [ -f /usr/bin/php-fpm7 ]; then
    echo "Using /usr/bin/php-fpm7"
    exec /usr/bin/php-fpm7 --nodaemonize
elif [ -f /usr/local/sbin/php-fpm ]; then
    echo "Using /usr/local/sbin/php-fpm"
    exec /usr/local/sbin/php-fpm --nodaemonize
elif [ -f /usr/local/bin/php-fpm ]; then
    echo "Using /usr/local/bin/php-fpm"
    exec /usr/local/bin/php-fpm --nodaemonize
elif command -v php-fpm8 >/dev/null 2>&1; then
    echo "Using php-fpm8 from PATH"
    exec php-fpm8 --nodaemonize
elif command -v php-fpm7 >/dev/null 2>&1; then
    echo "Using php-fpm7 from PATH"
    exec php-fpm7 --nodaemonize
elif command -v php-fpm >/dev/null 2>&1; then
    echo "Using php-fpm from PATH"
    exec php-fpm --nodaemonize
else
    echo "No PHP-FPM executable found!"
    # As a last resort, try to create a simple PHP server
    if command -v php >/dev/null 2>&1; then
        echo "Attempting to start PHP built-in server as fallback"
        cd /var/www/html/public || exit
        exec php -S 0.0.0.0:9000
    else
        echo "No PHP executable found. Cannot start PHP service."
        sleep 30
        exit 1
    fi
fi
