# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Firefly III
# Runs the PHP-FPM daemon
# ==============================================================================

# Create php-fpm config with proper user settings
mkdir -p /etc/php82/php-fpm.d
cat > /etc/php82/php-fpm.d/www.conf << EOF
[global]
error_log = /proc/1/fd/2
log_level = warning

[www]
; Use nobody user with security.limit_extensions=off to avoid permission issues
user = nobody
group = nobody
listen = 127.0.0.1:9000
pm = dynamic
pm.max_children = 10
pm.start_servers = 3
pm.min_spare_servers = 2
pm.max_spare_servers = 5
pm.max_requests = 500

; Security settings - disabling restrictions to make it work in container
security.limit_extensions = 
php_admin_flag[allow_url_fopen] = on

; Redirect logs to stdout/stderr
access.log = /proc/1/fd/1
php_admin_value[error_log] = /proc/1/fd/2
php_admin_flag[log_errors] = on
php_admin_value[memory_limit] = 512M
php_admin_value[upload_max_filesize] = 50M
php_admin_value[post_max_size] = 50M
php_admin_value[max_execution_time] = 300
php_admin_flag[display_errors] = off
php_admin_flag[expose_php] = off

; Fix for the invalid boolean value error
php_admin_value[cgi.fix_pathinfo] = 0
EOF

# Create required directories
directories=("/var/www/html/storage/logs" 
             "/var/www/html/storage/app/public"
             "/var/www/html/storage/framework/cache"
             "/var/www/html/storage/framework/sessions"
             "/var/www/html/storage/framework/views"
             "/var/www/html/bootstrap/cache")

# Ensure directories exist with permissive permissions
for dir in "${directories[@]}"; do
    mkdir -p "$dir" || true
    chmod -R 777 "$dir" || true
done

# Ensure Laravel log file exists with permissive permissions
touch /var/www/html/storage/logs/laravel.log || true
chmod 666 /var/www/html/storage/logs/laravel.log || true

# Also ensure the whole application directory is world-writable
chmod -R 777 /var/www/html || true

# Create the PHP-FPM tmp directories with permissive permissions
mkdir -p /tmp/php-fpm || true
chmod -R 777 /tmp/php-fpm || true

# Start PHP-FPM with explicit config file
bashio::log.info "Starting PHP-FPM..."
# Use --allow-to-run-as-root flag to override security restriction
exec php-fpm82 --nodaemonize --fpm-config /etc/php82/php-fpm.d/www.conf --allow-to-run-as-root
