# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Firefly III
# Runs the PHP-FPM daemon
# ==============================================================================

# Create php-fpm config with proper user settings
mkdir -p /etc/php82/php-fpm.d
cat > /etc/php82/php-fpm.d/www.conf << EOF
[global]
error_log = /proc/self/fd/2
log_level = warning

[www]
; Run as root to avoid permission issues
user = root
group = root
listen = 127.0.0.1:9000
pm = dynamic
pm.max_children = 10
pm.start_servers = 3
pm.min_spare_servers = 2
pm.max_spare_servers = 5
pm.max_requests = 500

; Redirect logs to stdout/stderr
access.log = /proc/self/fd/1
php_admin_value[error_log] = /proc/self/fd/2
php_admin_flag[log_errors] = on
php_admin_value[memory_limit] = 512M
php_admin_value[upload_max_filesize] = 50M
php_admin_value[post_max_size] = 50M
php_admin_value[max_execution_time] = 300
php_admin_flag[display_errors] = off
php_admin_flag[expose_php] = off
EOF

# Start PHP-FPM
bashio::log.info "Starting PHP-FPM..."
exec php-fpm82 --nodaemonize --fpm-config /etc/php82/php-fpm.d/www.conf
