# rootfs/etc/services.d/firefly/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Firefly III
# Take down the S6 supervision tree when queue worker fails
# ==============================================================================

if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "Firefly III queue worker crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Firefly III queue worker stopped, restarting..."
