# rootfs/etc/services.d/firefly/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Firefly III
# Runs the Firefly III queue worker
# ==============================================================================

# Wait until the database is available
bashio::net.wait_for "$(bashio::config 'database.host')" "$(bashio::config 'database.port')"

# Get the ingress URL from the Supervisor
INGRESS_URL=$(bashio::addon.ingress_url)
bashio::log.info "Firefly III is accessible at: ${INGRESS_URL}"

# Run the queue worker (artisan command) in the background
cd /var/www/html || exit

# Run the queue worker with proper error handling
exec php artisan queue:work
