# rootfs/etc/services.d/firefly/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Firefly III
# Runs the Firefly III queue worker
# ==============================================================================

cd /var/www/html || exit

# If queue:work fails, we'll just sleep indefinitely instead of crashing
if php artisan queue:work --sleep=3 --tries=3; then
  bashio::log.info "Queue worker started successfully"
else
  bashio::log.warning "Queue worker failed to start, entering sleep mode"
  # Just sleep indefinitely instead of crashing
  exec tail -f /dev/null
fi
